# supervisord.conf
[supervisord]
nodaemon=true

[program:plexivity]
command=python /plexivity/plexivity.py
